# etc/service/pd-server.service
[Unit]
Description=PD Server
Wants=network-online.target

[Service]
User=tidb
Group=tidb
ExecStart=/usr/bin/pd-server --config $CONFIGURATION_DIRECTORY/pd/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
KillSignal=SIGINT
ConfigurationDirectory=pd
WorkingDirectory=/var/lib/pd/

[Install]
WantedBy=multi-user.target
